# files/daphne.conf
[program:daphne]
command = /opt/app/start_daphne.sh
directory=/opt/app/dsbd
user = root
autostart=true
autorestart=true
stdout_logfile = /logs/daphne.log
redirect_stderr = true