# files/gunicorn.conf
[program:gunicorn]
command=/opt/app/start_gunicorn.sh
directory=/opt/app/dsbd
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile = /logs/gunicorn.log